# scripts/setup/foto-uart-updater.service
# Servicio systemd para auto-actualización de FotoUART
# Instalar con: sudo cp foto-uart-updater.service /etc/systemd/system/
# Habilitar con: sudo systemctl enable foto-uart-updater.timer

[Unit]
Description=FotoUART Auto-Updater
Documentation=https://github.com/cpatagon/foto-uart-dropin
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/home/pi/foto-uart-dropin
ExecStart=/home/pi/foto-uart-dropin/scripts/setup/auto_update.sh
StandardOutput=journal
StandardError=journal

# Seguridad
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/pi /var/log
ProtectHome=read-only